# deploy/eve-api.service
[Unit]
Description=EVE Online API Service
After=network.target

[Service]
Type=simple
User=eve-api
Group=eve-api
WorkingDirectory=/home/arturas/eve-api
ExecStart=java -jar /home/arturas/eve-api/eve-api.jar
Restart=always
RestartSec=10
Environment="SPRING_PROFILES_ACTIVE=prod"

[Install]
WantedBy=multi-user.target 